# etc/cron.d/resize
@reboot		root	/sbin/gpart resize -i mmcsd0

